# vuln-web-base/supervisord-mysqld.conf
[program:mysqld]
command=mysqld_safe
numprocs=1
autostart=true
autorestart=true