# sql-so-hard/supervisord-update-blacklist.conf
[program:nodejs]
command=/update_blacklist.sh
numprocs=1
autostart=true
autorestart=true
user=root
group=root