# sql-so-hard/supervisord-nodejs-app.conf
[program:nodejs]
command=/usr/bin/nodejs /app/app.js
numprocs=1
autostart=true
autorestart=true
user=root
group=root